# etc/systemd/system/gpu_exporter.service
[Unit]
Description=Export Intel GPU metrics to Node Exporter

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gpu_exporter.sh
